# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --touch-reload=/home/docker/code/uwsgi.ini --ini /home/docker/code/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx